# contester-advfiler.socket
[Unit]
Description=Contester filer
 
[Socket]
ListenStream=9094
 
[Install]
WantedBy=sockets.target
